# deploy/systemd/line-calendar-watch.service.template
[Unit]
Description=LINE Tsukishima calendar one-shot watch
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__WORKDIR__
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=__WORKDIR__/.env
ExecStart=__WORKDIR__/venv/bin/python __WORKDIR__/watch_calendar.py
User=__USER__
Group=__USER__

[Install]
WantedBy=multi-user.target



